# build/shop-api/data/conf/etc/supervisord.d/flask-api.ini
[program:flask-api]
command=/root/.pyenv/versions/api/bin/gunicorn --timeout 300  -w 4 --worker-class=gevent --max-requests 300 -b0.0.0.0:8010 manage:app
# command=command=uwsgi --ini /app/drcc-backend/uwsgi.ini
directory=/app/api/flask
user=root
startsecs=0
stopwaitsecs=60
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=
