# build/shop-api/data/conf/etc/supervisord.d/gin-api.ini
[program:gin-api]
command=/app/api/gin/ginserver
directory=/app/api/gin
user=root
startsecs=0
stopwaitsecs=60
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=
