# build/shop-api/data/conf/etc/supervisord.d/shop-celery.ini
[program:shop-celery]
command=/root/.pyenv/versions/api/bin/celery -A worker.celery worker --loglevel INFO
directory=/app/api/flask
user=root
startsecs=0
stopwaitsecs=60
stopsignal=QUIT
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
environment=
